# deploy/gmo-trader.service
[Unit]
Description=GMO Coin Trading Bot
After=network.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/opt/bitdaytrader
Environment="PATH=/opt/bitdaytrader/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/bitdaytrader"
ExecStart=/opt/bitdaytrader/.venv/bin/python -m src.main
Restart=always
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/gmo-trader/trading.log
StandardError=append:/var/log/gmo-trader/trading-error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/bitdaytrader/data /opt/bitdaytrader/logs /opt/bitdaytrader/models

# Resource limits
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
